# source/supervisord.ini
[program:tasks]
command=/usr/bin/php __project_dir__/think tasks prod
directory=__project_dir__
autostart=true
autorestart=true
startsecs=5
user=php-running userName
numprocs=64
redirect_stderr         = true
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups  = 10
stdout_logfile          = __project_dir__/runtime/log/task.log